# files/root/etc/supervisor/conf.d/cloud9.conf
[program:cloud9]
command=pm2 start /etc/pm2_cloud9.json
